# cdn-proxy.service
[Unit]
Description=CDN Proxy — Wasabi cache via Go Fiber
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/etc/cdn-proxy
ExecStart=/usr/local/bin/cdn-proxy
EnvironmentFile=/etc/cdn-proxy/.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
